# launchd/com.m13v.social-autoposter-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/matthewdi</string>
		<key>PATH</key>
		<string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
		<key>AUTOPOSTER_API_BASE</key>
		<string>https://s4l.ai</string>
	</dict>
	<key>Label</key>
	<string>com.m13v.social-autoposter-heartbeat</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/matthewdi/social-autoposter/scripts/heartbeat.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>900</integer>
	<key>StandardErrorPath</key>
	<string>/Users/matthewdi/social-autoposter/skill/logs/launchd-heartbeat-stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/matthewdi/social-autoposter/skill/logs/launchd-heartbeat-stdout.log</string>
</dict>
</plist>
