# launchd/com.m13v.social-gsc-seo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.m13v.social-gsc-seo</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/matthewdi/social-autoposter/seo/cron_gsc.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>StandardOutPath</key>
    <string>/Users/matthewdi/social-autoposter/seo/logs/gsc_launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/matthewdi/social-autoposter/seo/logs/gsc_launchd_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Library/Frameworks/Python.framework/Versions/3.12/bin</string>
        <key>HOME</key>
        <string>/Users/matthewdi</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
