# launchd/com.m13v.social-audit-linkedin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.m13v.social-audit-linkedin</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/matthewdi/social-autoposter/skill/audit-linkedin.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>4</integer>
		<key>Minute</key>
		<integer>15</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/matthewdi/social-autoposter/skill/logs/launchd-audit-linkedin-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/matthewdi/social-autoposter/skill/logs/launchd-audit-linkedin-stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/matthewdi/.nvm/versions/node/v20.19.4/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
		<key>HOME</key>
		<string>/Users/matthewdi</string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
