# launchd/com.m13v.social-reddit-threads.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/matthewdi</string>
		<key>PATH</key>
		<string>/Users/matthewdi/.nvm/versions/node/v20.19.4/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
	</dict>
	<key>Label</key>
	<string>com.m13v.social-reddit-threads</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/matthewdi/social-autoposter/skill/run-reddit-threads.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>/Users/matthewdi/social-autoposter/skill/logs/launchd-reddit-threads-stderr.log</string>
	<key>StandardOutPath</key>
	<string>/Users/matthewdi/social-autoposter/skill/logs/launchd-reddit-threads-stdout.log</string>
	<!--
	  4 fires per local day, aligned with twitter-threads (8/12/17/22) so the
	  dashboard "Post Threads" row shows one unified cadence across platforms.
	  Offset to :13 to avoid colliding with twitter-threads (:47).
	-->
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Hour</key>
			<integer>8</integer>
			<key>Minute</key>
			<integer>13</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>12</integer>
			<key>Minute</key>
			<integer>13</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>17</integer>
			<key>Minute</key>
			<integer>13</integer>
		</dict>
		<dict>
			<key>Hour</key>
			<integer>22</integer>
			<key>Minute</key>
			<integer>13</integer>
		</dict>
	</array>
</dict>
</plist>
