# launchd/com.m13v.social-seo-top-pages.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.m13v.social-seo-top-pages</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/matthewdi/social-autoposter/seo/run_top_pages_pipeline.sh</string>
    </array>
    <!-- Daily at 03:30 local (after seo-improve at 03:00); iterates every
         project with landing_pages.top_pages_enabled = true in config.json.
         Ranks pages by (pageviews + email_signups*100 + schedule_clicks*500
         + get_started_clicks*300 + bookings*1000) over 24h and generates
         an adjacent landing page riding the winner's topical momentum. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/matthewdi/social-autoposter/skill/logs/launchd-seo-top-pages-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/matthewdi/social-autoposter/skill/logs/launchd-seo-top-pages-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/matthewdi/.nvm/versions/node/v20.19.4/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/Library/Frameworks/Python.framework/Versions/3.12/bin</string>
        <key>HOME</key>
        <string>/Users/matthewdi</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
